# packages/ttio-web/units/ttio-web.target
[Unit]
Description=timetab.io Web
After=docker.service
Requires=docker.service
Wants=ttio-api.service \
      ttio-frontend.service \
      ttio-proxy.service \
      ttio-survey.service \
      ttio-workers.target

[Install]
WantedBy=default.target
