# packages/ttio-web/units/ttio-web.target
[Unit]
Description=timetab.io Web
After=docker.service
Requires=docker.service
Wants=ttio-proxy.service

[Install]
WantedBy=default.target
