# packages/ttio-server/cron.d/renew-certs
0 0 1 * * root /usr/local/bin/renew-certs
